# au.com.phdesign.500px-osx-background.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
	<dict>
		<key>Label</key>
		<string>au.com.phdesign.500px-osx-background</string>
		<key>ProgramArguments</key>
		<array>
			<string>/bin/sh</string>
			<string>/Users/paul/Projects/500px-osx-background/500px-osx-background.sh</string>
			<string>/Users/paul/Pictures/Wallpaper</string>
		</array>
		<key>EnvironmentVariables</key>
		<dict>
			<key>PATH</key>
			<string>/bin:/usr/bin:/usr/local/bin</string>
		</dict>
		<key>StartCalendarInterval</key>
		<dict>
			<key>Hour</key>
			<integer>12</integer>
			<key>Minute</key>
			<integer>0</integer>
		</dict>
		<key>StandardOutPath</key>
		<string>/tmp/500px-osx-background.log</string>
		<key>StandardErrorPath</key>
		<string>/tmp/500px-osx-background.log</string>
	</dict>
</plist>
